# recipes-core/cute-files/cute-files/cute-files.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/cute-files

[Install]
WantedBy=multi-user.target
